# configs/systemduser/vdirsyncer.service
[Unit]
Description=Synchronize vdirs with vdirsyncer

[Service]
ExecStart=/usr/bin/vdirsyncer sync
